# deploy/neogrip.service
[Unit]
Description=neogrip, prothèse sécurisé
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eegsvc
Group=eegsvc
StateDirectory=neogrip
StateDirectoryMode=0755
WorkingDirectory=/var/lib/neogrip
EnvironmentFile=/etc/neogrip/secrets.env

ExecStart=/usr/bin/python3 /opt/neogrip/main.py

#Restart=on-failure
#RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target